# package/base-files/files/etc/init.d/gps_bd
#!/bin/sh /etc/rc.common
# /init.d/set_rate
START=99

start() {
	#exit 0
	ser2net -c /etc/gps_bd.conf
}

